# pla_node/deploy/pla-node.service
[Unit]
Description=HexForge PLA Node (edge tool-runner API)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/pla_node
EnvironmentFile=/etc/pla_node/pla.env
ExecStart=/opt/pla_node/.venv/bin/gunicorn --workers 2 --bind 0.0.0.0:8787 app.main:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
